# sysd-services/network_monitor.service
[Unit]
Description=Network Monitor Service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/network_monitor.sh
Restart=always

[Install]
WantedBy=multi-user.target
